# solidity-platform-force.service
[Unit]
Description=Solidity Learning Platform (Force Deployment)
After=network.target

[Service]
Type=simple
User=elo
WorkingDirectory=/home/elo/learning_solidity
ExecStart=/usr/bin/npm run start:force
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

# Increase limits for error handling
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
